# factorio.path
[Unit]
Description=Starts the factorio server when its been downloaded and has its options environment file.

[Path]
Unit=factorio.service
PathExists=/var/cache/factoriod/factorio/bin/x64/factorio
PathExists=/var/cache/factoriod/factorio.opts.env

[Install]
WantedBy=multi-user.target
